# deploy/systemd/omnisight-backend.service
; systemd unit for the OmniSight FastAPI backend.
;
; Install:
;   sudo cp deploy/systemd/omnisight-backend.service /etc/systemd/system/
;   # Edit the USER_HOME / USERNAME placeholders below to match your WSL user.
;   sudo systemctl daemon-reload
;   sudo systemctl enable --now omnisight-backend
;
; I10: multi-worker mode. Set OMNISIGHT_WORKERS in .env (default: CPU/2, min 2).
; Requires OMNISIGHT_REDIS_URL for cross-worker shared state + SSE delivery.
;
; Graceful shutdown: the FastAPI lifespan cancels background loops per worker
; (watchdog / decision sweep / DLQ retry / IQ nightly / finetune nightly /
; memory decay / pubsub listener). KillSignal=SIGTERM triggers backend/lifecycle.py
; drain coordinator (G1 #1): stop accepting new requests, flush SSE, close DB,
; wait up to 30s for in-flight tasks. TimeoutStopSec=40 = 30s drain + 10s buffer
; before systemd escalates to SIGKILL.

[Unit]
Description=OmniSight backend (FastAPI / uvicorn)
Documentation=https://github.com/your-org/OmniSight-Productizer
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=exec
User=USERNAME
WorkingDirectory=/home/USER_HOME/work/sora/OmniSight-Productizer
EnvironmentFile=/home/USER_HOME/work/sora/OmniSight-Productizer/.env
ExecStart=/bin/sh -c '\
    WORKERS=$${OMNISIGHT_WORKERS:-$$(python3 -c "import os; print(max(2, os.cpu_count() // 2))")}; \
    exec python3 -m uvicorn backend.main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers $$WORKERS \
    --no-access-log'
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=40

; Sandboxing. OmniSight itself needs docker + /home access; keep the
; rest locked down so a compromise doesn't walk out.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/USER_HOME/work/sora/OmniSight-Productizer/data
ReadWritePaths=/home/USER_HOME/work/sora/OmniSight-Productizer/workspaces
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
